# services/pgbackrest/run
#!/command/execlineb -P
with-contenv
s6-setuidgid pgbackrest
pgbackrest server
