# services.d/cron/run
#!/command/execlineb -P

with-contenv

/usr/sbin/cron -f
