# services.d/tail/run
#!/command/execlineb -P

fdclose 2
/usr/bin/tail -F /var/log/pgbackrest/pgbackrest.log
